# mes/systemd/mes_swing_live.service
[Unit]
Description=MES Swing Trader (LIVE)
After=network.target

[Service]
Type=oneshot
User=ubu
Group=ubu
WorkingDirectory=/opt/mes

# Ensure HOME is set for op + Python
Environment=HOME=/home/ubu

# Safety arming
Environment="MES_SWING_ARMED=YES"

# 🔐 1Password service account token
EnvironmentFile=/etc/op.env

# LIVE secrets + environment mapping
ExecStart=/usr/bin/op run \
  --env-file=/opt/mes/env_live.op.env \
  -- /opt/mes/env_use_live.sh \
  /opt/mes/venv/bin/python \
  /opt/mes/mes_swing.py --allow-live

# Scheduling & safety (matches scalper expectations)
Nice=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

# 🔐 Allow 1Password to write its config only
ProtectHome=read-only
ReadWritePaths=/home/ubu/.config/op

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.targetLI
