# mes/systemd/mes_scalp_live.service
[Unit]
Description=MES Scalp Trader LIVE (30 min cycle — simplified runtime)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubu
WorkingDirectory=/opt/mes

# Environment may come from /etc/environment or config.json
EnvironmentFile=-/etc/op.env

# LIVE mode explicitly armed via flag + env
ExecStart=/usr/bin/op run \
  --env-file=/opt/mes/env_live.op.env \
  -- /opt/mes/env_use_live.sh \
  -- /opt/mes/venv/bin/python \
  /opt/mes/mes_scalp.py --allow-live --risk 0.02


StandardOutput=journal
StandardError=journal

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
