# leo/leo_server.timer
[Unit]
Description=Run Leo Server Telegram Report hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
