# leo/leo_server.timer
[Unit]
Description=Run Leo Server Telegram Report 6 hours

[Timer]
OnCalendar=*-*-* 0/6:00:00
Persistent=true

[Install]
WantedBy=timers.target