# kraken/kraken_sol.service
[Unit]
Description=Kraken SOL Trader – MES Crypto Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubu
WorkingDirectory=/home/ubu/leo-services/kraken

EnvironmentFile=/home/ubu/leo-services/secrets/global.op.env
EnvironmentFile=/home/ubu/leo-services/secrets/kraken_trade.op.env

ExecStart=/usr/bin/op run -- /usr/bin/python3 /home/ubu/leo-services/kraken/kraken_sol.py

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
